# cmd/relayer/install/systemd.service
[Unit]
Description=Relayer Service
After=network.target

[Service]
Type=simple
User=relayer
WorkingDirectory=/home/relayer
ExecStart=/usr/bin/relayer run --config /etc/relayer/config.yaml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
